# debian/nuntium.conf
description "nuntium service binds ofono and telepathy-ofono for MMS"

start on (:sys:started JOB=ofono or started dbus)
stop on :sys:stopping JOB=ofono

respawn 
respawn limit 10 10

exec nuntium
