# debian/nuntium.conf
description	"nuntium service binds ofono and telepathy-ofono for MMS"

start on :sys:started JOB=ofono
stop on :sys:stopped JOB=ofono

exec nuntium
